# root/etc/services.d/dash/run
#!/usr/bin/with-contenv bash

cd /app/app/server

export LINUX_DASH_SERVER_PORT=7171

exec \
	s6-setuidgid abc node index.js